# Ubuntu/wall_display/wm.target
[Unit]
Description=A Window Manager is running